# nut-2.7.2/scripts/systemd/nut-driver.service.in
[Unit]
Description=Network UPS Tools - power device driver controller
After=local-fs.target network.target udev-settle.service
Wants=udev-settle.service
StopWhenUnneeded=yes

[Service]
ExecStart=@SBINDIR@/upsdrvctl start
ExecStop=@SBINDIR@/upsdrvctl stop
Type=forking

