# lvm-common-1.5.17/debian/postrm
#!/bin/sh

set -e

if [ "$1" = "remove" -a -x /usr/sbin/update-devfsd ] ; then
	/usr/sbin/update-devfsd
fi

if [ "$1" = "purge" ] ; then
	rm -f /etc/lvmtab
	update-rc.d lvm remove > /dev/null
fi

exit 0

