# ntp-4.1.0/debian/ntp-simple.conffiles
/etc/cron.daily/ntp-simple
/etc/cron.weekly/ntp-simple
/etc/init.d/ntp
