# ntp-4.1.0/debian/ntp-refclock.conffiles
/etc/cron.daily/ntp-refclock
/etc/cron.weekly/ntp-refclock
/etc/init.d/ntp
