# apache2-2.4.25/debian/apache2.prerm
#!/bin/sh
set -e

# If upgrading, postinst will restart. So we only need stop if removing.
if [ -x "/etc/init.d/apache-htcacheclean" ] && [ "$1" = remove ] ; then
	invoke-rc.d apache-htcacheclean stop || true
fi

#DEBHELPER#
