# openldap2-2.0.23/debian/slapd.postrm
#!/bin/sh

set -e

purge_debconf() {
	if [ -e /usr/share/debconf/confmodule ] ; then
		. /usr/share/debconf/confmodule
		db_purge
	fi
}


purge_init() {
	update-rc.d slapd remove >/dev/null
}


if [ "$1" = "purge" ] ; then
	purge_init
	purge_debconf "$@"
	rm -f /etc/ldap/slapd.conf
	rm -rf /var/lib/ldap
fi

exit 0

