# openldap2-2.0.23/debian/slapd.conffiles
/etc/init.d/slapd
/etc/ldap/schema/corba.schema
/etc/ldap/schema/core.schema
/etc/ldap/schema/cosine.schema
/etc/ldap/schema/inetorgperson.schema
/etc/ldap/schema/java.schema
/etc/ldap/schema/krb5-kdc.schema
/etc/ldap/schema/misc.schema
/etc/ldap/schema/nis.schema
/etc/ldap/schema/openldap.schema
