# cups-1.7.5/data/cups.socket.in
[Unit]
Description=CUPS Printing Service Sockets

[Socket]
ListenStream=@CUPS_DEFAULT_DOMAINSOCKET@
BindIPv6Only=ipv6-only
PassCredentials=true

[Install]
WantedBy=sockets.target
