# cups-1.7.5/data/cups.path.in
[Unit]
Description=CUPS Printer Service Spool

[Path]
PathExistsGlob=@CUPS_REQUESTS@/d*

[Install]
WantedBy=paths.target
