# quota-4.00~pre1/debian/prerm
#! /bin/sh
# prerm script for quota
#
# see: dh_installdeb(1)

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see /usr/doc/packaging-manual/

set -e

case "$1" in
    remove|deconfigure)
	# debhelper only handles the quotarpc init script
	if [ -x "/etc/init.d/quota" ]; then
		if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
			invoke-rc.d quota stop || exit $?
		else
			/etc/init.d/quota stop || exit $?
		fi
	fi

	[ -f /var/lib/quota/off ] && rm -f /var/lib/quota/off
	[ -f /var/lib/quota/new ] && rm -f /var/lib/quota/new
        ;;
    upgrade)
	if [ -L /usr/doc/quota ]; then
	        rm -f /usr/doc/quota
	fi
        ;;
    failed-upgrade)
	if [ -L /usr/doc/quota ]; then
	        rm -f /usr/doc/quota
	fi
	;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

#DEBHELPER#	

exit 0
