# noflushd-2.8/debian/noflushd.service
[Unit]
Description=Idle Disk Spindown Daemon
After=udev-settle.service
Wants=udev-settle.service

[Service]
Environment=TIMEOUT=5
EnvironmentFile=/etc/default/noflushd
ExecStart=/usr/sbin/noflushd -d -n $TIMEOUT

[Install]
WantedBy=multi-user.target
