# raidtools2-0.90.20010914/debian/postrm
#! /bin/sh

set -e

. /usr/share/debconf/confmodule

if [ "$1" = "purge" ]
then
    update-rc.d raid2 remove > /dev/null
    test -f /etc/rcS.d/S??raid2 && rm -f /etc/rcS.d/S??raid2 || /bin/true
fi

if [ "$1" = purge -a -e /usr/share/debconf/confmodule ]; then
        . /usr/share/debconf/confmodule
        db_purge
fi
