# exim-3.35/debian/crontab
# /etc/cron.d/exim: crontab fragment for exim

# Run queue every 15 minutes
08,23,38,53 *     * * *     mail   if [ -x /usr/sbin/exim -a -f /etc/exim/exim.conf ]; then /usr/sbin/exim -q ; fi
