# cupsys-1.1.14/debian/cupsys.conffiles
/etc/cups/cupsd.conf
/etc/cups/mime.convs
/etc/cups/mime.types
/etc/init.d/cupsys
/etc/pam.d/cupsys
