# minidlna-1.1.6/debian/minidlna.service
[Unit]
Description=DLNA/UPnP-AV media server

[Service]
ExecStart=/usr/sbin/minidlnad -S

[Install]
WantedBy=multi-user.target
